# roles/ombi/templates/ombi.service.j2
[Unit]
Description=Systemd script to run Ombi as a service
After=network-online.target

[Service]
User=ombi
Group=ombi
Type=simple
ExecStart=/usr/bin/mono /opt/Ombi/Release/Ombi.exe

[Install]
WantedBy=multi-user.target
